# etc/crontab
# m h  dom mon dow   command

#Clear log
*/1 * * * * sh log.sh

#Backup Clear
0 2 */3 * * sh bckp.sh

#Mata Processo
*/3 * * * * killall -y 3m /usr/local/mkauth/backup/update
